# client/systemd/comConD.service
[Unit]
Description=ComConD

[Service]
ExecStart=/usr/bin/comConD.py

[Install]
WantedBy=multi-user.target
