# client/systemd/pmc_client.service
[Unit]
Description=pmc_client_updater

[Service]
ExecStart=/usr/bin/pmc_client HOSTNAME

[Install]
WantedBy=multi-user.target
