# client/systemd/pmcClient.service
[Unit]
Description=ComConD

[Service]
ExecStart=/usr/local/bin/pmcClient.py /etc/pmcClient.conf

[Install]
WantedBy=multi-user.target
